# scripts/ble-watchdog.service
[Unit]
Description=Recover hung onboard Bluetooth adapter (hci0)
ConditionPathExists=/usr/local/bin/ble-watchdog.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ble-watchdog.sh
